# cw-borg-backup.timer
[Unit]
Description=Run cw-borg-backup.service daily.

[Timer]
OnCalendar=07:00
AccuracySec=30m
Persistent=yes

[Install]
WantedBy=timers.target
